# systemd_user/pcloud.service
[Unit]
Description=Start PCloud
Wants=network-online.target

[Service]
Type=exec
ExecStart=/usr/bin/pcloud 


